# oss-mondobackup-full.service
[Unit]
Description=oss-mondobackup full backup
ConditionACPower=true

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/oss-mondobackup
ExecStart=/sbin/rcoss-mondobackup -f
